# rc.local
#!/bin/bash

sleep 60
/usr/bin/python /home/pi/Rpi/pifinder.py &
. /home/pi/Rpi/scripts/start_sensors.sh

/usr/bin/python3 /home/pi/Rpi/collect_adc.py &

exit 0
